# meta-evb/meta-evb-nuvoton/meta-evb-npcm845/recipes-phosphor/console/obmc-console/obmc-console-ttyS4-ssh.socket
[Unit]
Description=Phosphor Host Console SSH Per-Connection socket
Wants=obmc-console@ttyS4.service

[Socket]
ListenStream=8294
Accept=yes

[Install]
WantedBy=sockets.target